# nginx/controltower-api.service
# Control Tower API — systemd
# Copiar: sudo cp controltower-api.service /etc/systemd/system/
# Activar: sudo systemctl daemon-reload && sudo systemctl enable controltower-api && sudo systemctl start controltower-api
# Estado:  sudo systemctl status controltower-api
#
# Ajusta WorkingDirectory y ExecStart si tu repo está en otra ruta o usas otro usuario.

[Unit]
Description=Control Tower API (uvicorn)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/controltower/backend
Environment=BACKEND_HOST=127.0.0.1
Environment=BACKEND_PORT=8000
Environment=PATH=/var/www/controltower/backend/venv/bin:/usr/bin
ExecStart=/var/www/controltower/backend/venv/bin/python run_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
